# deploy/gpt-account-manager-upgrade-agent.timer
[Unit]
Description=Run GPT Account Manager upgrade agent every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=gpt-account-manager-upgrade-agent.service

[Install]
WantedBy=timers.target
